# safety-monitor/scripts/collector.service
[Unit]
Description=Safety Monitor Collector
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/safety-monitor/collector
ExecStart=/usr/bin/python3 /home/pi/safety-monitor/collector/main.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
