# safety-monitor/scripts/collector.service
[Unit]
Description=Safety Monitor Collector
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Siteguard/safety-monitor
ExecStart=/home/pi/.venv/bin/python -m collector.main
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target