# safety-monitor/scripts/ui.service
[Unit]
Description=Safety Monitor UI
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Siteguard/safety-monitor
ExecStart=/usr/bin/python3 -m uvicorn server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
